# .crontab
#!/usr/local/lib/python2.7/

46 09 * * * sh ~/git/sortes/sortes.sh
